# home-menu.service
[Unit]
Description=Home menu — project index page
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=simon
WorkingDirectory=/home/simon/projects
ExecStart=/usr/bin/python3 /home/simon/projects/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
